# rusty_ide_v2/scripts/rusty-agent.service
[Unit]
Description=Rusty IDE Agent Wrapper - File-based AI agent communication daemon
Documentation=https://github.com/yourusername/rusty-ide
After=network.target

[Service]
Type=simple
# Update this path to match your installation
ExecStart=/path/to/rusty_ide_v2/src-tauri/src/agent_wrapper.sh --daemon
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="RUSTY_AGENT_DIR=%h/.rusty/agent"
Environment="RUST_LOG=info"

# Security hardening (optional)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
